# InitApp_Task/meta-bbappendLayer/recipes-example/hellodaemon/files/hellodaemonAPP.service
[Unit]
Description=Script that runs hellodeamon

[Service]
Type=simple
ExecStart=/usr/bin/hellodeamon

[Install]
WantedBy=multi-user.target
